# etc/supervisor/conf.d/mbserverd.conf
[program:mbserverd]
command=mbserverd -n -p 5020
user=root
autostart=true
autorestart=false
