# etc/supervisor/conf.d/mbserverd.conf
[program:mbserverd]
command=mbserverd -n -p 502
user=root
autostart=true
autorestart=false
